# server/overlays/etc/supervisor.d/gunicorn.conf
[program:gunicorn]
command=/bin/sh -c ". /etc/lava-server/lava-server-gunicorn && /usr/bin/gunicorn3 lava_server.wsgi --log-level $LOGLEVEL --log-file $LOGFILE --workers $WORKERS"
autostart=false
