# server/overlays/etc/supervisor.d/postinst.conf
[program:postinst]
command=/usr/local/bin/python3 /usr/share/lava-server/postinst.py
autorestart=false
